# ansible/tmpl/rc.local.j2
#!/bin/sh -x

iwconfig wlan0 power off
exit 0
